# python_can_udp/DEBIAN/prerm
#!/bin/sh
# prerm script for sysmon

set -e

systemctl stop python_can_udp.service
systemctl disable python_can_udp.service

rm -rf /opt/python_can_udp/__pycache__

exit 0
